# go-backend/deploy/ark-lmd-go.service
[Unit]
Description=Ark LMD Go Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ark-lmd/go-backend
ExecStart=/opt/ark-lmd/go-backend/ark-lmd-go-backend
Restart=always
RestartSec=3

Environment=APP_ENV=production
Environment=PORT=3003
Environment=CORS_ORIGIN=https://ark-lmd.top
Environment=TRUST_PROXY=true
Environment=LOG_DIR=/var/log/ark-lmd
Environment=LOG_JSON=true
Environment=LOG_TO_STDOUT=false
Environment=LOG_LEVEL=info
Environment=LOG_IP_HASH=false
Environment=LOG_IP_HASH_SALT=replace-with-a-random-salt
Environment=MAX_CONCURRENCY=2
Environment=MAX_QUEUE_SIZE=4
Environment=CALC_TIMEOUT_MS=15000
Environment=RATE_LIMIT_PER_MINUTE=15
Environment=CACHE_TTL_SECONDS=3600
Environment=CACHE_MAX_ENTRIES=1024
Environment=DATA_FILE=/opt/ark-lmd/data/gameItems.json
Environment=ADMIN_TOKEN=replace-with-a-long-random-token

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/ark-lmd

[Install]
WantedBy=multi-user.target
